# scripts/hugin.service
[Unit]
Description=Hugin Network Scanner - Installation Validation
Documentation=https://github.com/Smarttfoxx/Hugin
After=network.target

[Service]
Type=oneshot
User=root
ExecStartPre=/bin/mkdir -p /var/log/hugin /var/lib/hugin
ExecStart=/bin/bash -c 'echo "Hugin Network Scanner is installed and ready" && /usr/local/bin/hugin --version'
RemainAfterExit=yes
TimeoutStartSec=30

# Environment
Environment=PATH=/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
Alias=hugin-scanner.service
